# ark-flask.service
[Unit]
Description=Project Ark — Flask captive portal (RAG over offline Wikipedia + Ollama)
After=network-online.target ark-kiwix.service ollama.service
Wants=network-online.target ark-kiwix.service ollama.service

[Service]
Type=simple
User=@INSTALL_USER@
WorkingDirectory=@ARK_DIR@
Environment=ARK_OLLAMA_MODEL=@OLLAMA_MODEL@
Environment=ARK_OLLAMA_URL=http://127.0.0.1:11434
Environment=ARK_KIWIX_URL=http://127.0.0.1:8080
Environment=ARK_PORT=80
Environment=ARK_DATA_DIR=@ARK_DATA_DIR@
Environment=PYTHONUNBUFFERED=1
ExecStart=@VENV_DIR@/bin/python @ARK_DIR@/app.py
Restart=on-failure
RestartSec=5

# Allow binding to privileged port 80 without running as root.
# (install.sh also runs setcap on the venv python as a belt-and-braces.)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
