# ark-llama-cpp.service
[Unit]
Description=Project Ark — llama.cpp inference server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@INSTALL_USER@
WorkingDirectory=/opt/llama-cpp
ExecStart=/opt/llama-cpp/server -m @MODELS_DIR@/@OLLAMA_MODEL@.gguf -ngl 999 -t 4 --port 8001
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax=6G
CPUQuota=80%

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/llama-cpp @MODELS_DIR@

[Install]
WantedBy=multi-user.target
